# launchd/co.rockstarsunlimited.starshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>co.rockstarsunlimited.starshot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-lc</string>
    <string>cd /path/to/starshot &amp;&amp; /path/to/bun run upload</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Users/YOUR_USER/Desktop</string>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/co.rockstarsunlimited.starshot.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/co.rockstarsunlimited.starshot.err</string>
</dict>
</plist>
